# .scripts/init_services/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin blas --noclear %I 38400 linux
Type=simple
