# .scripts/init_services/pacman-update.service
[Service]
ExecStart=/bin/bash /root/pacmanupdate

[Install]
WantedBy=multi-user.target

